# mcmrelay.service.example
[Unit]
Description=MCMRelay - MeshCore Matrix Bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/path/to/mcmrelay/.venv/bin/mcmrelay
WorkingDirectory=/path/to/mcmrelay
Restart=on-failure
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
